# systemd/invoiceninja-scheduler.timer
[Unit]
Description=Run Invoice Ninja scheduler every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
